# beanstalkd/supervisord.conf
[supervisord]
nodaemon=true

[program:beanstalkd]
command=beanstalkd
autostart=true
autorestart=true
